# crates/codexswitch-cli/systemd/codexswitch.service.d/10-maintenance-resources.conf
[Service]
Nice=10
CPUWeight=25
IOWeight=25
IOSchedulingClass=idle
MemoryHigh=4G
MemoryMax=6G
MemorySwapMax=2G
TasksMax=512
LimitNOFILE=65536
